# wpa_supplicant/.footprint
drwxr-xr-x	niewod/users	etc/
drwxr-xr-x	niewod/users	etc/rc.d/
-rwxr-xr-x	niewod/users	etc/rc.d/wpa_supplicant
-rw-------	niewod/users	etc/wpa_supplicant.conf
drwxr-xr-x	niewod/users	usr/
drwxr-xr-x	niewod/users	usr/etc/
drwxr-xr-x	niewod/users	usr/etc/dbus-1/
drwxr-xr-x	niewod/users	usr/etc/dbus-1/system.d/
-rw-r--r--	niewod/users	usr/etc/dbus-1/system.d/wpa_supplicant.conf
drwxr-xr-x	niewod/users	usr/man/
drwxr-xr-x	niewod/users	usr/man/man5/
-rw-r--r--	niewod/users	usr/man/man5/wpa_supplicant.conf.5.gz
drwxr-xr-x	niewod/users	usr/man/man8/
-rw-r--r--	niewod/users	usr/man/man8/wpa_background.8.gz
-rw-r--r--	niewod/users	usr/man/man8/wpa_cli.8.gz
-rw-r--r--	niewod/users	usr/man/man8/wpa_passphrase.8.gz
-rw-r--r--	niewod/users	usr/man/man8/wpa_supplicant.8.gz
drwxr-xr-x	niewod/users	usr/sbin/
-rwxr-xr-x	niewod/users	usr/sbin/wpa_cli
-rwxr-xr-x	niewod/users	usr/sbin/wpa_passphrase
-rwxr-xr-x	niewod/users	usr/sbin/wpa_supplicant
drwxr-xr-x	niewod/users	usr/share/
drwxr-xr-x	niewod/users	usr/share/dbus-1/
drwxr-xr-x	niewod/users	usr/share/dbus-1/system-services/
-rw-r--r--	niewod/users	usr/share/dbus-1/system-services/fi.w1.wpa_supplicant1.service
